# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command = /usr/local/bin/uwsgi --ini /app/conf/uwsgi.ini
